# clawbridge/rootfs/etc/services.d/clawbridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ClawBridge
# ==============================================================================

bashio::log.info "Starting ClawBridge..."

exec python3 /app/server.py
